# DEBIAN/postinst
systemctl enable xxx.service
systemctl start xxx.service
