# etc/systemd/system/xxx.service
[Unit]
Description=XXX

[Service]
Type=simple
WorkingDirectory=/tmp
ExecStart=/usr/bin/python -m SimpleHTTPServer

[Install]
WantedBy=default.target
